# deploy/systemd/bot6rangos-heartbeat.service
[Unit]
Description=bot6rangos heartbeat monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot6rangos
EnvironmentFile=/home/ubuntu/bot6rangos/.env
ExecStart=/home/ubuntu/bot6rangos/.venv/bin/python /home/ubuntu/bot6rangos/heartbeat_monitor.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target
