# deploy/systemd/bot6rangos-telegram.service
[Unit]
Description=bot6rangos Telegram commands
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bot6rangos
EnvironmentFile=/home/ubuntu/bot6rangos/.env
ExecStart=/home/ubuntu/bot6rangos/.venv/bin/python /home/ubuntu/bot6rangos/telegram_bot_commands.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
